# configs/openrc-mailcloak
#!/sbin/openrc-run

name="mailcloak"
command="/usr/local/sbin/mailcloak"
command_args="/etc/mailcloak/config.yaml"
command_background="yes"
pidfile="/run/mailcloak.pid"

depend() {
  need net
  after postfix
}

start_pre() {
  checkpath -d -m 0750 -o root:root /etc/mailcloak
  checkpath -d -m 0750 -o mailcloak:dovecot /var/lib/mailcloak
}

stop_post() {
  rm -f /var/spool/postfix/private/mailcloak-policy /var/spool/postfix/private/mailcloak-socketmap
}